# pack/DEBIAN/postinst.in
# /bin/bash

systemctl enable @PROJECT_NAME@.service
systemctl daemon-reload
service @PROJECT_NAME@ start