# pack/DEBIAN/prerm.in
# /bin/bash

if systemctl --all --type service | grep -q "@PROJECT_NAME@";then
    systemctl disable @PROJECT_NAME@.service
    service @PROJECT_NAME@ stop
fi
